# files/conf/consul_template.ini
[program:consul_template]
command=consul-template -config=/etc/consul_t.conf
stopsignal=INT
user=consul
stdout_logfile=/var/log/consultemplate.log
autostart=true
autorestart=true
